# 99-serial_console.rules
KERNELS=="ttyUSB*", SUBSYSTEM=="tty", RUN+="/usr/local/bin/handle_serial_event"